# GopherServer.Service/systemdService/gopherServer.service
[Unit]
Description=Gopher Server

[Service]
Type=notify
ExecStart=/home/veso266/gopherServer/GopherServer.Service
Restart=on-failure
RestartSec=1s
StartLimitIntervalSec=0

[Install]
WantedBy=multi-user.target